# painel-ips.service
[Unit]
Description=Painel IPS
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/painel-ips
ExecStart=/usr/bin/env node /opt/painel-ips/server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production
Environment=PAINEL_HOST=0.0.0.0
Environment=PAINEL_PORT=3500

[Install]
WantedBy=multi-user.target
